# 99-predatortune.rules
SUBSYSTEM=="module", KERNEL=="predatortune_fan", RUN+="/bin/chmod 666 /sys/kernel/predatortune/fan_speed"
